# dot_systemwide/etc/cron.hourly/executable_btrbk
#!/bin/sh
mkdir -p /boot.backup
rsync -au --delete /boot/ /boot.backup/
exec /usr/bin/btrbk -q snapshot
