# private_dot_config/systemd/user/emacs.service
[Unit]
Description=Emacs daemon

[Service]
ExecStart=/usr/bin/emacs --daemon

[Install]
WantedBy=default.target
